# collis-family-planner.service
[Unit]
Description=Collis Family Planner Flask Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/collis-family-planner
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/doppler run -- /home/ubuntu/collis-family-planner/venv/bin/python app.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/ubuntu/collis-family-planner/logs/app.log
StandardError=append:/home/ubuntu/collis-family-planner/logs/app.log

[Install]
WantedBy=multi-user.target
